# packages/unisound-agent-frame/unisound_agent_frame-0.0.1.tar.gz/unisound_agent_frame-0.0.1/unisound_agent_frame/scripts/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB


[program:unisound-agent-web]
command=python service.py --env prod
directory=/app/unisound-agent-web
autostart=true
autorestart=true
stdout_logfile=/var/log/unisound-agent-web.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/unisound-agent-web.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
